# scripts/com.user.obsidian.dailycommit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.obsidian.dailycommit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/edwardbaeg/dev/dotfiles/scripts/daily-obsidian-commit.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/edwardbaeg/Library/Logs/obsidian-dailycommit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/edwardbaeg/Library/Logs/obsidian-dailycommit.error</string>
</dict>
</plist>